# platform/gentoo/files/etc/profile.d/makeopts.sh
#!/bin/bash
export MAKEOPTS="-j$(($(nproc) / 2))"
